# etc/hotplug.d/iface/81-speedtest
#!/bin/sh
[ "$ACTION" == "ifup" ] && [ "$INTERFACE" == "wan" ] && {
    /etc/init.d/auto_speedtest wan_up
}
