# etc/crontabs/root
*/3 * * * * /usr/sbin/ntpsetclock rand once >/dev/null 2>&1
5 */8 * * * /usr/sbin/sec_clt.sh 600 >/dev/null 2>&1
0 8,18 * * * /usr/sbin/crontab_record_camera.lua
* * * * * /usr/sbin/startscene_crontab.lua `/bin/date "+%u %H:%M"`
0 12 * * * /usr/sbin/recordscene_crontab.lua
0 5 */1 * * /usr/sbin/ad_filter.sh 1800 >/dev/null 2>&1
45 23 * * * /usr/sbin/points_sysset_pro.lua >/dev/null 2>&1
0 20 * * * /usr/bin/stat_lan
5 * * * * /usr/sbin/otapredownload >/dev/null 2>&1
